# host-setup/container-network-rules.timer
[Unit]
Description=Apply container network rules shortly after boot

[Timer]
OnBootSec=3min

[Install]
WantedBy=timers.target
