# host-setup/container-network-rules.service
[Unit]
Description=Restore iptables rules for container access (TFTP/Unimus automation)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/container-network-rules.sh
